# update-motd.d/15-fastfetch
#!/bin/sh

export PATH="/home/linuxbrew/.linuxbrew/bin:$PATH"
fastfetch
