# update-motd.d/40-fastfetch
#!/bin/bash

# Try to locate fastfetch
FASTFETCH_PATH="$(command -v fastfetch)"

# If not found, try Homebrew default path
if [ -z "$FASTFETCH_PATH" ]; then
    if [ -x "/home/linuxbrew/.linuxbrew/bin/fastfetch" ]; then
        export PATH="/home/linuxbrew/.linuxbrew/bin:$PATH"
        FASTFETCH_PATH="/home/linuxbrew/.linuxbrew/bin/fastfetch"
    fi
fi

REAL_USER=$(id -nu 1000 2>/dev/null || echo "$SUDO_USER")
USER_HOME=$(getent passwd "$REAL_USER" | cut -d: -f6)
CONFIG_PATH="$USER_HOME/.config/fastfetch/config.jsonc"

# Run fastfetch if it's now available
if [ -n "$FASTFETCH_PATH" ]; then
    fastfetch --config "$CONFIG_PATH"
else
    echo "fastfetch not found"
fi

